# rootfs/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

echo "Starting Jackett."
umask 002
exec s6-setuidgid jackett sh -c "HOME=/config mono /opt/Jackett/JackettConsole.exe -d /config"
